# .config/production/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/backend/.config/production/uwsgi_socket.ini

[program:front]
command=cd front;ng serve

[program:nginx]
command=nginx
